# deploy/d-brain-weekly.service
[Unit]
Description=d-brain Weekly Digest
After=network.target

[Service]
Type=oneshot
User=dbrain
WorkingDirectory=/home/dbrain/agent-second-brain
ExecStart=/home/dbrain/.local/bin/uv run python scripts/weekly.py
EnvironmentFile=/home/dbrain/agent-second-brain/.env
StandardOutput=journal
StandardError=journal
