# deploy/d-brain-process.service
[Unit]
Description=d-brain Daily Processing

[Service]
Type=oneshot
User=dbrain
WorkingDirectory=/home/dbrain/agent-second-brain
ExecStart=/home/dbrain/agent-second-brain/scripts/process.sh
Environment=PYTHONUNBUFFERED=1
